# openipc/stop_ap_mode.service
[Unit]
Description="stop ap mode service"
After=multi-user.target
StartLimitIntervalSec=10
StartLimitBurst=0

[Service]
Type=simple
ExecStart=bash /config/scripts/stop_ap_mode.sh

[Install]
WantedBy=multi-user.target
